# sabnzbdvpn.conf
[program:start]
autorestart = false
startsecs = 0
user = root
command = /root/start.sh

[program:sabnzbd]
user = nobody
command = /home/nobody/sabnzbd.sh
